# systemd/hyprland.target
[Unit]
Description=Hyprland
After=graphical.target
AllowIsolate=yes
